# systemd/turtle-visualizer.service
[Unit]
Description=turtle-visualizer Hydra sketch server
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/turtle-visualizer
Environment=PORT=8080
ExecStart=/usr/bin/node /opt/turtle-visualizer/server.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
