# systemd/turtle-visualizer-kiosk.service
[Unit]
Description=turtle-visualizer kiosk browser
Wants=turtle-visualizer.service
After=turtle-visualizer.service

[Service]
Type=simple
User=pi
Environment=APP_URL=http://localhost:8080
Environment=HOME=/home/pi
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/turtle-visualizer/scripts/start-kiosk.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
